# packaging/root/etc/rc.d/init.d/rundeckd
#!/bin/sh
#
# rundeckd    Startup script for the tapd
#
# chkconfig: 2345 90 10
# description: rundeckd, providing rundeckd
# pidfile: /var/run/rundeckd.pid

# Source function library
. /etc/rc.d/init.d/functions
. /etc/rundeck/client/profile

prog="rundeckd"
rundeckd="${JAVA_HOME:-/usr}/bin/java ${RDECK_JVM} ${BOOTSTRAP_CP}:/etc/rundeck/server com.dtolabs.rundeck.RunServer /etc/rundeck/client 9999"
RETVAL=0
PID_FILE=/var/run/rundeckd.pid

start() {
	echo -n $"Starting $prog: "
	nohup $rundeckd 2>&1 >/var/log/rundeck/service.log &
	RETVAL=$?
	echo $RETVAL > /var/run/rundeckd.pid
	echo
	[ $RETVAL -eq 0 ] && touch /var/lock/subsys/$prog
	return $RETVAL
}

stop() {
	echo -n $"Stopping $prog: "
	killproc $rundeckd
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/$prog
	return $RETVAL
}

reload() {
	echo -n $"Reloading $prog: "
	killproc $rundeckd -HUP
	RETVAL=$?
	echo
	return $RETVAL
}

case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		start
		;;
	condrestart)
		if [ -f /var/lock/subsys/$prog ]; then
			stop
			start
		fi
		;;
	reload)
		reload
		;;
	status)
		status $rundeckd
		RETVAL=$?
		;;
	*)
		echo $"Usage: $0 {start|stop|restart|condrestart|reload|status}"
		RETVAL=1
esac

exit $RETVAL
